# mkosi.images/initrd/mkosi.extra/usr/lib/systemd/system/systemd-repart.service.d/sysroot.conf
[Unit]
After=sysroot.mount
Before=initrd-switch-root.service
ConditionDirectoryNotEmpty=|/sysroot/usr/lib/repart.d
